# crontab
* * * * * node /app/cron.sh > /tmp/cronjob.log 2>&1
