# crontab
* * * * * node bash node queueFeeds.js
* * * * * node bash node updateFeeds.js
